# systemd/template.service.erb
[Unit]
Description=<%= description %>
<%- # This doesn't really do very much since systemd considers a
    # service "started" when it launches, and working around that for user
    # units is hard, so we actually handle the ordering in the
    # run_container code, but this at least helps somewhat
    # -%>
<%- if defined?(after_containers) -%>
Requires=<%= after_containers.split("\s").map { |x| "#{x}.service" }.join(' ') %>
After=<%= after_containers.split("\s").map { |x| "#{x}.service" }.join(' ') %>
<%- end -%>

[Service]
Restart=always
<%- if defined?(needs_network) -%>
ExecStartPre=/bin/bash -c 'while ! ping -c1 google.com; do sleep 1; done'
<%- end -%>
SyslogIdentifier=<%= bundle %>-<%= name %>
WorkingDirectory=<%= maindir %>
<%- if type == 'container' -%>
ExecStart=/bin/bash -x <%= maindir %>/run_container.sh <%= container_name %>
ExecStop=/bin/bash -x <%= maindir %>/destroy_container.sh <%= container_name %>
<%- end -%>
<%- if type == 'addon' -%>
ExecStart=/bin/bash -x <%= maindir %>/run_addon.sh <%= container_name %> <%= addon_name %>
ExecStop=/bin/bash -x <%= maindir %>/stop_addon.sh <%= container_name %> <%= addon_name %>
<%- end -%>
TimeoutStartSec=10min

[Install]
WantedBy=default.target
